# slotmachine_app/supervisord.conf
[supervisord]
nodaemon=true

[program:new_user]
command=python new_user.py
autostart=true
autorestart=true

[program:spin]
command=python spin.py
autostart=true
autorestart=true